# deploy/systemd/HomeMonitoringDashboard.service.d/10-wait-migration.conf
[Unit]
# Requires (not Wants) so the dashboard does NOT start if migrations fail — it must never
# run against an unmigrated/incompatible schema. After= sequences it behind the one-shot.
After=HomeMonitoringMigration.service
Requires=HomeMonitoringMigration.service
