# deploy/systemd/HomeMonitoringDashboard.service.d/10-wait-migration.conf
[Unit]
After=HomeMonitoringMigration.service
Wants=HomeMonitoringMigration.service
