# deploy/systemd/HomeMonitoringDashboard.service.d/20-environment.conf
[Service]
# HTTP listen address for the dashboard. This is a host/deployment concern, so it lives in the
# unit (as a drop-in) rather than appsettings.json: a deploy of the base appsettings.json can
# never wipe it, and it never affects Aspire dev (systemd isn't used there). 0.0.0.0 = listen on
# all interfaces (LAN-reachable); use 127.0.0.1 to expose the app only through a local proxy.
Environment=ASPNETCORE_URLS=http://0.0.0.0:5000
# Optional: restore OpenTelemetry export by pointing at your OTLP collector. The exporter reads
# this from the environment (not appsettings), so it belongs here too. Uncomment and set your host:
#Environment=OTEL_EXPORTER_OTLP_ENDPOINT=http://YOUR-OTLP-HOST:4317
